# changebg.service
[Unit]
Description=My Custom Task

[Service]
ExecStart=/home/root/backgrounds/changebg.sh

